# deploy/supervisor/worker-queue.conf
[program:horizon]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan horizon
autostart=true
autorestart=true
numprocs=1
startsecs = 0
startretries=3
stdout_events_enabled=1
loglevel=warn
redirect_stderr=true
user=www-data
